# setup/kittyhack_control.service
[Unit]
Description=KittyHack Control Service (remote control target)
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/root/kittyhack
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/root/kittyhack/.venv/bin"
ExecStart=/root/kittyhack/.venv/bin/python -m src.kittyhack_control
Restart=always
RestartSec=2

# Process handling
KillSignal=SIGTERM
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
